# .config/systemd/user/dictate-ptt.service
[Unit]
Description=Dictate push-to-talk client
After=dictate.service
Requires=dictate.service

[Service]
ExecStart=%h/.local/bin/dictate --ptt --paste
Restart=on-failure
RestartSec=5

[Install]
WantedBy=default.target
